# init-solr.conf
[program:init-solr]
user = solr
directory = /opt/solr
command = /bin/bash bin/init-solr.sh
autorestart = false