# mongo-connector.conf
[program:mongo-connector]
user = solr
directory = /opt/solr
command = /bin/bash bin/mongo-connector.sh
autorestart = false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0