# solr.conf
[program:solr]
user = solr
directory = /opt/solr
command = /bin/bash bin/start-solr.sh
autorestart = true
