# solr.conf
[program:solr]
user = solr
directory = /opt/solr
command = /opt/solr/bin/solr -f
autorestart = true
